# bts.conf
[supervisord]
logfile = /tmp/supervisord.log

[program:btsserver]
command=$GOPATH/bin/main
autostart=true
autorestart=true
redirect_stderr=true